# src/supervisor.d/dummie.ini
[program:dummie]
command=/bin/touch -c  c
autorestart=false
#startretries=0

