# src/supervisor.d/dummie.ini
[program:dummie]
command=
autorestart=false
startretries=0

