# src/supervisor.d/dummie.ini
[program:dummie]
command=/usr/bin/touch --no-create c
autorestart=false
#startretries=0

