# src/supervisor.d/dummie.ini
[program:dummie]
command=date
#autorestart=false
#startretries=0

