# src/supervisor.d/Fix-permission.ini
[program:Fix_acl]
#command=setfacl -m g:10001:rwx test/
command=/f_acl.sh
autorestart=false
startretries=0

