# src/supervisor.d/Fix-permission.ini
[program:Fix_acl]
command=/f_acl.sh
autorestart=false
startretries=0
autostart=true
